# DroneAdHoc/AdHoc.service
[Unit]
Description=SetUpAdHoc script

[Service]
ExecStart=/usr/bin/AdHocScript

[Install]
WantedBy=multi-user.target

